# android-service@hwcomposer.service.d/20-phosh.conf
[Unit]
BindsTo=phosh.service
Before=phosh.service
